# etc/s6-overlay/s6-rc.d/libseccomp2/up
foreground { echo "Checking libseccomp2" }
/etc/s6-overlay/scripts/libseccomp2_check.sh