# etc/s6-overlay/s6-rc.d/libseccomp2/up
#!/bin/sh
/etc/s6-overlay/scripts/libseccomp2_check.sh